# systemd/uiosync.timer
[Unit]
Description=Run uio pull sync every 5 min

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
