# systemd/uiosyncDown.service
[Unit]
Description=Pulls files from the remote server

[Service]
ExecStart=/home/sully/scripts/uiosync-down.sh

[Install]
WantedBy=default.target
